# cron/cron.hourly/run_hourly
#!/bin/sh

python3 /home/neil/nestock/instock/job/basic_data_daily_job.py
#mkdir -p /data/logs
#DATE=`date +%Y-%m-%d:%H:%M:%S`
#echo $DATE >> /data/logs/hourly.log